# zzz.conf
[program:zzz]
command=python3 /home/curator2/zzz/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/zzz/test.err.log
stdout_logfile=/home/curator2/zzz/test.out.log